# build/MountTab.service
[Unit]
Description=A system to manage tabs
After=network.target

[Service]
Type=simple
ExecStart=/home/tylord/.local/bin/mt serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
